# deploy/systemd/trading-research-daily.timer
[Unit]
Description=Run Trading Research Lab daily deterministic validation

[Timer]
OnCalendar=*-*-* 02:30:00 UTC
Persistent=true
Unit=trading-research-daily.service

[Install]
WantedBy=timers.target

