# deploy/systemd/trading-research-self-improvement.timer
[Unit]
Description=Run Trading Research Lab self-improvement cycle

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true
Unit=trading-research-self-improvement.service

[Install]
WantedBy=timers.target

