# deploy/systemd/trading-research-hourly.timer
[Unit]
Description=Run Trading Research Lab hourly source scan

[Timer]
OnCalendar=hourly
Persistent=true
Unit=trading-research-hourly.service

[Install]
WantedBy=timers.target

