# deploy/systemd/trading-research-weekly.service
[Unit]
Description=Trading Research Lab weekly deep research
After=network-online.target

[Service]
Type=oneshot
User=trading
Group=trading
WorkingDirectory=/opt/trading/research-lab
EnvironmentFile=-/opt/trading/research-lab/.env
ExecStartPre=/usr/bin/mkdir -p /opt/trading/research-lab/tmp
ExecStart=/usr/bin/flock -n /opt/trading/research-lab/tmp/research.lock /opt/trading/research-lab/.venv/bin/python /opt/trading/research-lab/scripts/run_weekly_deep_research.py
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/trading/research-lab
Restart=no
TimeoutStartSec=1200
